# TP 3/yt.service
[Unit]
Description=Youtube downloader

[Service]
ExecStart=sudo bash /home/elbatista/srv/yt/yt-v2.sh

[Install]
WantedBy=multi-user.target